# etc/periodic/daily/120.clean-preserve
#!/bin/sh
#
# $FreeBSD: src/etc/periodic/daily/120.clean-preserve,v 1.4 1999/08/27 23:24:03 peter Exp $
#

if [ -d /var/preserve ]; then
	echo ""
	echo "Removing stale files from /var/preserve:"

	cd /var/preserve && find . ! -name . -mtime +7 -delete
fi
