# etc/periodic/daily/120.clean-preserve
#!/bin/sh
#
# $FreeBSD$
#
# Remove stale files in /var/preserve
#

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

case "$daily_clean_preserve_enable" in
    [Yy][Ee][Ss])
	if [ -n "$daily_clean_preserve_days" -a -d /var/preserve ]
	then
	    echo ""
	    echo "Removing stale files from /var/preserve:"

	    case "$daily_clean_preserve_verbose" in
		[Yy][Ee][Ss])
		    print=-print;;
		*)
		    print=;;
	    esac

	    cd /var/preserve &&
		find . ! -name . -mtime +$daily_clean_preserve_days \
		    -delete $print
	fi;;
esac
