# etc/periodic/daily/130.clean-msgs
#!/bin/sh
#
# $Id: 130.clean-msgs,v 1.1.1.1 1997/08/16 17:04:00 pst Exp $
#
# remove system	messages older than 21 days
#
if [ -d /var/msgs ] ; then
    echo ""
    echo "Cleaning out old system announcements:"

    msgs -c
fi
