# etc/periodic/daily/140.clean-rwho
#!/bin/sh
#
# $FreeBSD: src/etc/periodic/daily/140.clean-rwho,v 1.4 1999/08/27 23:24:03 peter Exp $
#

if [ -d /var/rwho ] ; then
	echo ""
	echo "Removing stale files from /var/rwho:"

	cd /var/rwho && find . ! -name . -mtime +7 -delete
fi

