# etc/periodic/daily/140.clean-rwho
#!/bin/sh
#
# $Id: 140.clean-rwho,v 1.2 1997/08/17 17:55:30 pst Exp $
#

if [ -d /var/rwho ] ; then
	echo ""
	echo "Removing stale files from /var/rwho:"

	cd /var/rwho && find . ! -name . -mtime +7 -delete
fi

