# etc/periodic/daily/140.clean-rwho
#!/bin/sh
#
# $Id: 140.clean-rwho,v 1.1.1.1 1997/08/12 17:51:15 pst Exp $
#

if [ -d /var/rwho ] ; then
	echo ""
	echo "Removing stale files from /var/rwho"
	cd /var/rwho && find . ! -name . -mtime +7 -exec rm -f -- {} \;
fi

