# etc/periodic/daily/150.clean-hoststat
#!/bin/sh
#
# $FreeBSD$
#

if [ -d /var/spool/.hoststat ] ; then
	echo ""
	echo "Removing stale files from /var/spool/.hoststat:"

	cd /var/spool/.hoststat && find * -mtime +3 -delete
fi
