# etc/periodic/daily/150.clean-hoststat
#!/bin/sh
#
# $Id: 140.clean-rwho,v 1.3 1999/03/21 12:33:54 joerg Exp $
#

if [ -d /var/spool/.hoststat ] ; then
	echo ""
	echo "Removing stale files from /var/spool/.hoststat:"

	cd /var/rwho && find * -mtime +3 -delete
fi
