# etc/periodic/daily/150.clean-hoststat
#!/bin/sh
#
# $Id: 150.clean-hoststat,v 1.1 1999/03/30 15:30:44 wollman Exp $
#

if [ -d /var/spool/.hoststat ] ; then
	echo ""
	echo "Removing stale files from /var/spool/.hoststat:"

	cd /var/spool/.hoststat && find * -mtime +3 -delete
fi
