# etc/periodic/daily/320.rdist
#!/bin/sh
#
# $FreeBSD: src/etc/periodic/daily/320.rdist,v 1.2.2.1 1999/08/29 14:19:40 peter Exp $
#

if [ -f /etc/Distfile ]; then
    echo ""
    echo "Running rdist with /etc/Distfile:"

    rdist -f /etc/Distfile
fi
