# etc/periodic/daily/320.rdist
#!/bin/sh
#
# $Id: 320.rdist,v 1.2 1997/08/17 17:55:32 pst Exp $
#

if [ -f /etc/Distfile ]; then
    echo ""
    echo "Running rdist with /etc/Distfile:"

    rdist -f /etc/Distfile
fi
