# etc/periodic/daily/320.rdist
#!/bin/sh
#
# $Id: 320.rdist,v 1.1.1.1 1997/08/16 17:04:00 pst Exp $
#

if [ -f /etc/Distfile ]; then
    echo ""
    echo "Running rdist with /etc/Distfile:"

    rdist -f /etc/Distfile
fi
