# etc/periodic/daily/340.uucp
#!/bin/sh
#
# $Id: 340.uucp,v 1.1.1.1 1997/08/16 17:04:01 pst Exp $
#
# Local cleanup of UUCP files.  This is for backwards compatibility,
# /etc/uuclean.daily doesn't exist by default.
#
if [ -d /var/spool/uucp -a -f /etc/uuclean.daily ]; then
	echo ""
	echo "Cleaning up UUCP:"

	echo /etc/uuclean.daily | su daemon
fi

