# etc/periodic/daily/400.status-disks
#!/bin/sh
#
# $FreeBSD: src/etc/periodic/daily/400.status-disks,v 1.1.1.1.2.1 1999/08/29 14:19:41 peter Exp $
#

echo ""
echo "Disk status:"

df -k -t nonfs

# display which filesystems need backing up

echo ""
dump W
