# etc/periodic/daily/400.status-disks
#!/bin/sh
#
# $FreeBSD$
#

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

case "$daily_status_disks_enable" in
    [Yy][Ee][Ss])
	echo ""
	echo "Disk status:"

	df $daily_status_disks_df_flags && rc=0 || rc=3

	# display which filesystems need backing up

	echo ""
	dump W || rc=3;;

    *)  rc=0;;
esac

exit $rc
