# etc/periodic/daily/400.status-disks
#!/bin/sh
#
# $FreeBSD: src/etc/periodic/daily/400.status-disks,v 1.2 1999/08/27 23:24:04 peter Exp $
#

echo ""
echo "Disk status:"

df -k -t nonfs

# display which filesystems need backing up

echo ""
dump W
