# etc/periodic/daily/400.status-disks
#!/bin/sh
#
# $FreeBSD$
#

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

case "$daily_status_disks_enable" in
    [Yy][Ee][Ss])
	echo ""
	echo "Disk status:"

	df -k -t nonfs

	# display which filesystems need backing up

	echo ""
	dump W;;
esac
