# etc/periodic/daily/410.status-uucp
#!/bin/sh
#
# $Id: 410.status-uucp,v 1.1.1.1 1997/08/12 17:51:15 pst Exp $
#

if [ -d /var/spool/uucp -a -x /usr/bin/uustat ]; then
    echo ""
    echo "uucp status:"
    uustat -a
fi
