# etc/periodic/daily/420.status-network
#!/bin/sh
#
# $FreeBSD: src/etc/periodic/daily/420.status-network,v 1.2.2.1 1999/08/29 14:19:42 peter Exp $
#

if [ -x /usr/bin/netstat ] ; then
    echo ""
    echo "Network interface status:"

    netstat -i
fi
