# etc/periodic/daily/430.status-rwho
#!/bin/sh
#
# $Id: 430.status-rwho,v 1.2 1998/10/06 09:46:46 joerg Exp $
#

if [ -d /var/rwho -a -x /usr/bin/rwho -a $(ls -l /var/rwho | wc -l) -ne 0 ]
then
    echo ""
    echo "Local network system status:"
    ruptime
else
    echo ""
    echo "Local system status:"
    uptime
fi
