# etc/periodic/daily/430.status-rwho
#!/bin/sh
#
# $FreeBSD$
#

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

case "$daily_status_rwho_enable" in
    [Yy][Ee][Ss])
	rwho=$(echo /var/rwho/*)
        if [ -x /usr/bin/rwho -a -f "${rwho%% *}" ]
        then
	    echo ""
	    echo "Local network system status:"
	    ruptime
	else
	    echo ""
	    echo "Local system status:"
	    uptime
	fi;;
esac
