# etc/periodic/daily/450.status-security
#!/bin/sh
#
# $FreeBSD: src/etc/periodic/daily/450.status-security,v 1.2.2.1 1999/08/29 14:19:43 peter Exp $
#

if [ -f /etc/security -a -x /usr/sbin/sendmail ] ; then
    echo ""
    echo "Security check:"
    echo "    (output mailed separately)"

    sh /etc/security 2>&1 | sendmail root
fi
