# etc/periodic/daily/450.status-security
#!/bin/sh
#
# $FreeBSD$
#

if [ -f /etc/security -a -x /usr/sbin/sendmail ] ; then
    echo ""
    echo "Security check:"
    echo "    (output mailed separately)"

    sh /etc/security 2>&1 | sendmail root
fi
