# etc/periodic/daily/450.status-security
#!/bin/sh
#
# $Id: 450.status-security,v 1.1.1.1 1997/08/16 17:04:01 pst Exp $
#

if [ -f /etc/security -a -x /usr/sbin/sendmail ] ; then
    echo ""
    echo "Security check:"
    echo "    (output mailed separately)"

    sh /etc/security 2>&1 | sendmail root
fi
