# etc/periodic/daily/460.status-mail-rejects
#!/bin/sh
#
# $FreeBSD$
#

if [ -d /etc/mail -a -f /var/log/maillog ]; then
  echo
  echo Checking for rejected mail hosts:

  start=`date -v-1d '+%b %d' | sed 's/0\(.\)$/ \1/'`
  zcat -fc /var/log/maillog.0* /var/log/maillog | grep reject= |
	perl -ne "print \"\$2\n\" 
	if /^$start.*ruleset=check_\S+,\s+arg1=(<[^@]+@)?([^>,]+).*reject=/o;" |
	sort | uniq -c | sort -nr
fi
