# etc/periodic/daily/480.leapfile-ntpd
#!/bin/sh
#
# $FreeBSD$
#

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

case "$daily_ntpd_leapfile_enable" in
    [Yy][Ee][Ss])
	case "$daily_ntpd_leapfile_background" in
	[Yy][Ee][Ss])
	    (anticongestion && service ntpd onefetch) &
	    ;;
	*)
	    anticongestion
	    service ntpd onefetch
	    ;;
	esac
	;;
esac

exit $rc
