# etc/periodic/daily/999.local
#!/bin/sh
#
# $Id: 999.local,v 1.1.1.1 1997/08/12 17:51:16 pst Exp $
#
# Run the old /etc/daily.local script.  This is really for backwards
# compatibility more than anything else.
#
if [ -f /etc/daily.local ]; then
    echo ""
    echo "Running daily.local:"

    sh /etc/daily.local
fi
