# etc/periodic/daily/Makefile
# $FreeBSD$

FILES=	100.clean-disks \
	110.clean-tmps \
	120.clean-preserve \
	130.clean-msgs \
	140.clean-rwho \
	150.clean-hoststat \
	200.backup-passwd \
	210.backup-aliases \
	300.calendar \
	310.accounting \
	330.news \
	400.status-disks \
	405.status-ata-raid \
	406.status-gmirror \
	420.status-network \
	430.status-rwho \
	440.status-mailq \
	450.status-security \
	460.status-mail-rejects \
	470.status-named \
	500.queuerun \
	999.local

.include <bsd.prog.mk>
