# rc.d/ypbind
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ypbind
# REQUIRE: ypserv
# BEFORE:  DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="ypbind"
command="/usr/sbin/${name}"
start_precmd="ypbind_precmd"

load_rc_config $name
rcvar="nis_client_enable"
command_args="${nis_client_flags}"

ypbind_precmd()
{
	local _domain

	if ! checkyesno rpcbind_enable  && \
	    ! /etc/rc.d/rpcbind forcestatus 1>/dev/null 2>&1
	then
		force_depend rpcbind || return 1
	fi

	_domain=`domainname`
	if [ -z "$_domain" ]; then
		warn "NIS domainname(1) is not set."
		return 1
	fi
}

run_rc_command "$1"
