# etc/rc.d/abi
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: abi
# REQUIRE: archdep
# KEYWORD: nojail

. /etc/rc.subr

name="abi"
start_cmd="${name}_start"
stop_cmd=":"

sysv_start()
{
	echo -n ' sysvipc'
	load_kld sysvmsg
	load_kld sysvsem
	load_kld sysvshm
}

linux_start()
{
	echo -n ' linux'
	load_kld -e 'linux(aout|elf)' linux
	if [ -x /compat/linux/sbin/ldconfig ]; then
		/compat/linux/sbin/ldconfig
	fi
}

svr4_start()
{
	echo -n ' svr4'
	load_kld -m svr4elf svr4
}

abi_start()
{
	echo -n 'Additional ABI support:'

	checkyesno sysvipc_enable && sysv_start
	checkyesno linux_enable && linux_start
	checkyesno svr4_enable && svr4_start

	echo '.'
}

load_rc_config $name
run_rc_command "$1"
