# etc/rc.d/abi
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: abi
# REQUIRE: LOGIN
# BEFORE:  securelevel
# KEYWORD: nojail

. /etc/rc.subr

echo -n 'Additional ABI support:'

name="sysvipc"
rcvar=`set_rcvar`
start_cmd="sysv_start"
stop_cmd=":"

sysv_start()
{
	echo -n ' sysvipc'
	kldload sysvmsg >/dev/null 2>&1
	kldload sysvsem >/dev/null 2>&1
	kldload sysvshm >/dev/null 2>&1
}
load_rc_config $name
run_rc_command "$1"

name="linux"
rcvar=`set_rcvar`
start_cmd="linux_start"

linux_start()
{
	echo -n ' linux'
	if ! kldstat -v | grep -E 'linux(aout|elf)' > /dev/null; then
		kldload linux > /dev/null 2>&1
	fi
	if [ -x /compat/linux/sbin/ldconfig ]; then
		/compat/linux/sbin/ldconfig
	fi
}
load_rc_config $name
run_rc_command "$1"

name="svr4"
rcvar=`set_rcvar`
start_precmd="echo -n ' svr4'"
start_cmd="kldload svr4 > /dev/null 2>&1"

load_rc_config $name
run_rc_command "$1"

echo '.'
