# etc/rc.d/adjkerntz
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/adjkerntz,v 1.4.2.1 2004/10/10 09:50:53 mtm Exp $
#

# PROVIDE: adjkerntz
# REQUIRE: mountcritlocal random
# BEFORE: netif
# KEYWORD: nojail

. /etc/rc.subr

name="adjkerntz"
start_cmd="adjkerntz -i"
stop_cmd=":"

load_rc_config $name
run_rc_command "$1"
