# etc/rc.d/adjkerntz
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/adjkerntz,v 1.6.10.1.8.1 2012/03/03 06:15:13 kensmith Exp $
#

# PROVIDE: adjkerntz
# REQUIRE: FILESYSTEMS random
# BEFORE: netif
# KEYWORD: nojail

. /etc/rc.subr

name="adjkerntz"
start_cmd="adjkerntz -i"
stop_cmd=":"

load_rc_config $name
run_rc_command "$1"
