# etc/rc.d/adjkerntz
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/adjkerntz,v 1.5.2.1 2007/05/24 16:14:37 des Exp $
#

# PROVIDE: adjkerntz
# REQUIRE: FILESYSTEMS random
# BEFORE: netif
# KEYWORD: nojail

. /etc/rc.subr

name="adjkerntz"
start_cmd="adjkerntz -i"
stop_cmd=":"

load_rc_config $name
run_rc_command "$1"
