# etc/rc.d/adjkerntz
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/adjkerntz,v 1.5.2.1.4.1 2008/10/02 02:57:24 kensmith Exp $
#

# PROVIDE: adjkerntz
# REQUIRE: FILESYSTEMS random
# BEFORE: netif
# KEYWORD: nojail

. /etc/rc.subr

name="adjkerntz"
start_cmd="adjkerntz -i"
stop_cmd=":"

load_rc_config $name
run_rc_command "$1"
