# etc/rc.d/altqd
#!/bin/sh
#
# $NetBSD: altqd,v 1.2 2002/03/22 04:33:58 thorpej Exp $
# $FreeBSD$
#

# PROVIDE: altqd
# REQUIRE: network mountcritlocal dhclient
# BEFORE:  NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="altqd"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/altq.conf"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
