# etc/rc.d/amd
#!/bin/sh
#
# $NetBSD: amd,v 1.9 2002/03/22 04:33:58 thorpej Exp $
#

# PROVIDE: amd
# REQUIRE: rpcbind mountall ypbind
# BEFORE:  DAEMON

. /etc/rc.subr

name="amd"
rcvar=$name
command="/usr/sbin/${name}"
load_rc_config $name

command_args='-p -a '$amd_dir' -F /etc/amd.conf >/var/run/amd.pid'
required_files="/etc/amd.conf"
required_dirs="$amd_dir"
required_vars="rpcbind"

run_rc_command "$1"
