# etc/rc.d/apmd
#!/bin/sh
#
# $NetBSD: apmd,v 1.4 2000/09/19 13:04:38 lukem Exp $
#

# PROVIDE: apmd
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="apmd"
rcvar=$name
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
