# etc/rc.d/apmd
#!/bin/sh
#
# $NetBSD: apmd,v 1.3 2000/05/13 08:45:06 lukem Exp $
#

# PROVIDE: apmd
# REQUIRE: DAEMON

. /etc/rc.subr

name="apmd"
rcvar=$name
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
