# etc/rc.d/apmd
#!/bin/sh
#
# $NetBSD: apmd,v 1.5 2002/03/22 04:33:58 thorpej Exp $
# $FreeBSD$
#

# PROVIDE: apmd
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="apmd"
rcvar=$name
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
