# etc/rc.d/auto_linklocal
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: auto_linklocal
# REQUIRE: root
# BEFORE:  sysctl

. /etc/rc.subr

name="auto_linklocal"
rcvar=`set_rcvar ipv6`
start_cmd="${SYSCTL_W} net.inet6.ip6.auto_linklocal=1"
stop_cmd=":"

load_rc_config $name
run_rc_command "$1"
