# etc/rc.d/bgfsck
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: bgfsck
# REQUIRE: cron devfs
# KEYWORD: FreeBSD

. /etc/rc.subr

name="background-fsck"
rcvar="background_fsck"
start_precmd="echo 'Starting background file system checks.'"
start_cmd="nice -4 fsck -B -p 2>&1 | logger -p daemon.notice &"
stop_cmd=":"

load_rc_config $name
run_rc_command "$1"
