# etc/rc.d/bootparams
#!/bin/sh
#
# $NetBSD: bootparams,v 1.6 2002/03/22 04:33:58 thorpej Exp $
# $FreeBSD$
#

# PROVIDE: bootparams
# REQUIRE: rpcbind DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="bootparamd"
rcvar=$name
command="/usr/sbin/rpc.${name}"
required_files="/etc/bootparams"

load_rc_config $name
run_rc_command "$1"
