# etc/rc.d/bootparams
#!/bin/sh
#
# $NetBSD: bootparams,v 1.6 2002/03/22 04:33:58 thorpej Exp $
# $FreeBSD$
#

# PROVIDE: bootparams
# REQUIRE: rpcbind DAEMON
# BEFORE:  LOGIN
# KEYWORD: FreeBSD

. /etc/rc.subr

name="bootparamd"
rcvar=`set_rcvar`
required_files="/etc/bootparams"

case `${CMD_OSTYPE}` in
FreeBSD)
	command="/usr/sbin/${name}"
	;;
NetBSD)
	command="/usr/sbin/rpc.${name}"
	;;
esac

load_rc_config $name
run_rc_command "$1"
