# etc/rc.d/bootparams
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/bootparams,v 1.7.10.1.6.1 2010/02/10 00:26:20 kensmith Exp $
#

# PROVIDE: bootparams
# REQUIRE: rpcbind DAEMON
# BEFORE:  LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name="bootparamd"
rcvar=`set_rcvar`
required_files="/etc/bootparams"
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
