# etc/rc.d/bsnmpd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/bsnmpd,v 1.2.2.2.2.1 2012/03/03 06:15:13 kensmith Exp $
#

# PROVIDE: bsnmpd
# REQUIRE: NETWORKING syslogd
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="bsnmpd"
rcvar="bsnmpd_enable"
command="/usr/sbin/${name}"
pidfile="/var/run/snmpd.pid"

load_rc_config $name
run_rc_command "$1"
