# etc/rc.d/ccd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/ccd,v 1.5.2.1.2.1 2008/10/02 02:57:24 kensmith Exp $
#

# PROVIDE: disks
# KEYWORD: nojail

. /etc/rc.subr

name="ccd"
start_cmd="ccd_start"
stop_cmd=":"

ccd_start()
{
	if [ -f /etc/ccd.conf ]; then
		echo "Configuring CCD devices."
		ccdconfig -C
	fi
}

load_rc_config $name
run_rc_command "$1"
