# etc/rc.d/cleartmp
#!/bin/sh
#
# $NetBSD: cleartmp,v 1.4 2002/03/22 04:33:58 thorpej Exp $
# $FreeBSD$
#

# PROVIDE: cleartmp
# REQUIRE: mountall
# BEFORE: DAEMON
# KEYWORD: FreeBSD NetBSD

. /etc/rc.subr

name="cleartmp"
rcvar=`set_rcvar clear_tmp`
start_cmd="cleartmp_start"
stop_cmd=":"

cleartmp_start()
{
	echo "Clearing /tmp."
	#
	#	Prune quickly with one rm, then use find to clean up
	#	/tmp/[lq]* (this is not needed with mfs /tmp, but
	#	doesn't hurt anything).
	#
	(cd /tmp && rm -rf [a-km-pr-zA-Z]* &&
	    find -x . ! -name . ! -name lost+found ! -name quota.user \
		! -name quota.group -exec rm -rf -- {} \; -type d -prune)
	
	case ${OSTYPE} in
	FreeBSD)
		# Remove X lock files, since they will prevent you from
		# restarting X
		#
		rm -f /tmp/.X*-lock
		rm -fr /tmp/.X11-unix
		mkdir -m 1777 /tmp/.X11-unix
		;;
	NetBSD)
		;;
	esac
}

load_rc_config $name
run_rc_command "$1"
