# etc/rc.d/cleanvar
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: cleanvar
# REQUIRE: FILESYSTEMS var

. /etc/rc.subr

name="cleanvar"
rcvar=`set_rcvar`

start_precmd="${name}_prestart"
start_cmd="${name}_start"
stop_cmd=":"

extra_commands="reload"
reload_cmd="${name}_start"

purgedir()
{
	local dir file

	if [ $# -eq 0 ]; then
		purgedir .
	else
		for dir
		do
		(
			cd "$dir" && for file in .* *
			do
				# Skip over logging sockets
				[ -S "$file" -a "$file" = "log" ] && continue
				[ -S "$file" -a "$file" = "logpriv" ] && continue
				[ ."$file" = .. -o ."$file" = ... ] && continue
				if [ -d "$file" -a ! -L "$file" ]
				then
					purgedir "$file"
				else
					rm -f -- "$file"
				fi
			done
		)
		done
	fi
}

cleanvar_prestart()
{
	# These files must be removed only the first time this script is run
	# on boot.
	#
	rm -f /var/run/clean_var /var/spool/lock/clean_var
}

cleanvar_start () 
{
	if [ -d /var/run -a ! -f /var/run/clean_var ]; then
		purgedir /var/run
		# And an initial utmpx active session file
		(cd /var/run && cp /dev/null utx.active && chmod 644 utx.active)
		>/var/run/clean_var
	fi
	if [ -d /var/spool/lock -a ! -f /var/spool/lock/clean_var ]; then
		purgedir /var/spool/lock
		>/var/spool/lock/clean_var
	fi
	rm -rf /var/spool/uucp/.Temp/*
}

load_rc_config $name
run_rc_command "$1"
