# etc/rc.d/ddb
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/ddb,v 1.1.2.1.6.1 2010/02/10 00:26:20 kensmith Exp $
#

# PROVIDE: ddb
# REQUIRE: dumpon
# BEFORE: disks savecore initrandom
# KEYWORD: nojail

. /etc/rc.subr

name="ddb"
rcvar=`set_rcvar`
command="/sbin/${name}"
start_precmd="ddb_prestart"
stop_cmd=":"

ddb_prestart()
{
	# Silently exit if ddb is not enabled
	if [ -z "`sysctl -Nq debug.ddb.scripting.scripts`" ]; then
		return 1
	fi
}

load_rc_config $name

required_files="${ddb_config}"
command_args="${ddb_config}"

run_rc_command "$1"
