# etc/rc.d/devfs
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: devfs
# REQUIRE: LOGIN
# KEYWORD: FreeBSD

. /etc/rc.subr

name="devfs"

load_rc_config $name

# Setup DEVFS, ie permissions, links etc.
#
if [ -r /etc/rc.devfs ]; then
	sh /etc/rc.devfs
fi
