# etc/rc.d/dhcpd
#!/bin/sh
#
# $NetBSD: dhcpd,v 1.4 2000/09/19 13:04:38 lukem Exp $
#

# PROVIDE: dhcpd
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="dhcpd"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf /var/db/${name}.leases"

load_rc_config $name
run_rc_command "$1"
