# etc/rc.d/dhcpd
#!/bin/sh
#
# $NetBSD: dhcpd,v 1.3 2000/05/13 08:45:06 lukem Exp $
#

# PROVIDE: dhcpd
# REQUIRE: DAEMON

. /etc/rc.subr

name="dhcpd"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf /var/db/${name}.leases"

load_rc_config $name
run_rc_command "$1"
