# etc/rc.d/dmesg
#!/bin/sh
#
# $NetBSD: dmesg,v 1.8 2002/03/22 04:33:58 thorpej Exp $
# $FreeBSD$
#

# PROVIDE: dmesg
# REQUIRE: mountcritremote
# BEFORE:  DAEMON
# KEYWORD: Daemon nojail

. /etc/rc.subr

name="dmesg"
rcvar=`set_rcvar`
dmesg_file="/var/run/dmesg.boot"
start_cmd="do_dmesg"
stop_cmd=":"

do_dmesg()
{
	rm -f ${dmesg_file}
	( umask 022 ; /sbin/dmesg $rc_flags > ${dmesg_file} )
}

load_rc_config $name
run_rc_command "$1"
