# etc/rc.d/ftp-proxy
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/ftp-proxy,v 1.1.2.2.2.1 2009/04/15 03:14:26 kensmith Exp $
#

# PROVIDE: ftp-proxy
# REQUIRE: DAEMON pf
# KEYWORD: shutdown

. /etc/rc.subr

name="ftpproxy"
rcvar=`set_rcvar`
command="/usr/sbin/ftp-proxy"

load_rc_config $name
run_rc_command "$1"
