# etc/rc.d/ftp-proxy
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ftp-proxy
# REQUIRE: DAEMON pf

. /etc/rc.subr

name="ftpproxy"
rcvar=`set_rcvar`
command="/usr/sbin/ftp-proxy"

load_rc_config $name
run_rc_command "$1"
