# etc/rc.d/ftp-proxy
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/ftp-proxy,v 1.3.2.1.6.1 2010/12/21 17:09:25 kensmith Exp $
#

# PROVIDE: ftp-proxy
# REQUIRE: DAEMON pf
# KEYWORD: shutdown

. /etc/rc.subr

name="ftpproxy"
rcvar=`set_rcvar`
command="/usr/sbin/ftp-proxy"

load_rc_config $name
run_rc_command "$1"
