# etc/rc.d/gssd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/gssd,v 1.1.2.1.6.1 2010/12/21 17:09:25 kensmith Exp $
#

# PROVIDE: gssd
# REQUIRE: root
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="gssd"

load_rc_config $name
rcvar="gssd_enable"
command="${gssd:-/usr/sbin/${name}}"
eval ${name}_flags=\"${gssd_flags}\"
run_rc_command "$1"
