# etc/rc.d/hostapd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/hostapd,v 1.1.2.1 2006/02/12 22:13:15 flz Exp $
#

# PROVIDE: hostapd
# REQUIRE: mountcritremote
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="hostapd"
command="/usr/sbin/${name}"
rcvar=`set_rcvar`

conf_file="/etc/${name}.conf"

command_args="-B ${conf_file}"
required_files="${conf_file}"

load_rc_config ${name}
run_rc_command "$1"
