# etc/rc.d/idmapd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/idmapd,v 1.1.2.1.4.1 2010/02/10 00:26:20 kensmith Exp $
#

# PROVIDE: idmapd
# REQUIRE: rpcbind
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="idmapd"

load_rc_config $name
rcvar="idmapd_enable"
command="${idmapd:-/sbin/${name}}"
eval ${name}_flags=\"${idmapd_flags}\"
run_rc_command "$1"
