# etc/rc.d/ike
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/ike,v 1.4.2.1 2006/01/22 13:24:58 yar Exp $
#

# PROVIDE: ike
# REQUIRE: root mountcritlocal
# BEFORE:  DAEMON
# KEYWORD: nojail

. /etc/rc.subr

name="ike"
rcvar=`set_rcvar`
command="/usr/local/sbin/isakmpd"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
