# etc/rc.d/ike
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/ike,v 1.1 2003/11/28 17:28:42 bms Exp $
#

# PROVIDE: ike
# REQUIRE: root beforenetlkm mountcritlocal
# BEFORE:  DAEMON
# KEYWORD: FreeBSD

. /etc/rc.subr

name="ike"
rcvar=`set_rcvar`
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
