# etc/rc.d/inetd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/inetd,v 1.6.2.2.2.1 2012/03/03 06:15:13 kensmith Exp $
#

# PROVIDE: inetd
# REQUIRE: DAEMON LOGIN cleanvar
# KEYWORD: shutdown

. /etc/rc.subr

name="inetd"
rcvar="inetd_enable"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
