# etc/rc.d/initrandom
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: initrandom
# REQUIRE: dumpon ddb
# BEFORE: disks
# KEYWORD: nojail

. /etc/rc.subr

name="initrandom"
start_cmd="initrandom_start"
stop_cmd=":"

initrandom_start()
{
	soft_random_generator=`sysctl kern.random 2>/dev/null`

	echo -n 'Entropy harvesting:'

	if [ \! -z "${soft_random_generator}" ] ; then

		if [ -w /dev/random ]; then
			if checkyesno harvest_interrupt; then
				${SYSCTL} kern.random.sys.harvest.interrupt=1 >/dev/null
				echo -n ' interrupts'
			else
				${SYSCTL} kern.random.sys.harvest.interrupt=0 >/dev/null
			fi

			if checkyesno harvest_ethernet; then
				${SYSCTL} kern.random.sys.harvest.ethernet=1 >/dev/null
				echo -n ' ethernet'
			else
				${SYSCTL} kern.random.sys.harvest.ethernet=0 >/dev/null
			fi

			if checkyesno harvest_p_to_p; then
				${SYSCTL} kern.random.sys.harvest.point_to_point=1 >/dev/null
				echo -n ' point_to_point'
			else
				${SYSCTL} kern.random.sys.harvest.point_to_point=0 >/dev/null
			fi

			if checkyesno harvest_swi; then
				${SYSCTL} kern.random.sys.harvest.swi=1 >/dev/null
				echo -n ' swi'
			else
				${SYSCTL} kern.random.sys.harvest.swi=0 >/dev/null
			fi
		fi

	fi

	echo '.'
}

load_rc_config random
run_rc_command "$1"
