# etc/rc.d/ipmon
#!/bin/sh
#
# $NetBSD: ipmon,v 1.6 2002/03/22 04:33:58 thorpej Exp $
#

# PROVIDE: ipmon
# REQUIRE: syslogd
# BEFORE:  SERVERS

. /etc/rc.subr

name="ipmon"
rcvar=$name
command="/usr/sbin/${name}"
command_args="-D"

load_rc_config $name
run_rc_command "$1"
