# etc/rc.d/isdnd
#!/bin/sh
#
# $NetBSD: isdnd,v 1.9 2002/04/10 23:37:13 martin Exp $
# $FreeBSD$
#
# Mostly based on original script (/etc/rc.isdn) written by  Hellmuth Michaelis
#

# PROVIDE: isdnd
# REQUIRE: netif mountcritlocal cleanvar
# KEYWORD: nojail

. /etc/rc.subr

name="isdnd"
rcvar=`set_rcvar isdn`
pidfile="/var/run/${name}.pid"
command="/usr/sbin/isdnd"
start_cmd="isdnd_start"

isdnd_start()
{
	echo -n 'ISDN subsystem setup:'

	# Start isdnd
	#
	echo -n ' isdnd'
	case ${isdn_fsdev} in
	[Nn][Oo] | '')
		/usr/sbin/isdnd ${isdn_flags}
		;;
	*)
		# Change vidmode of ${isdn_fsdev}
		#
		case ${isdn_screenflags} in
		[Nn][Oo])
			;;
		*)
			/usr/sbin/vidcontrol < ${isdn_fsdev} > ${isdn_fsdev} 2>&1 ${isdn_screenflags}
			;;
		esac

		/usr/sbin/isdnd ${isdn_flags} -f -r ${isdn_fsdev} -t ${isdn_ttype}
		;;
	esac

	# Start isdntrace
	#
	if checkyesno isdn_trace; then
		echo -n ' isdntrace'
		nohup /usr/sbin/isdntrace ${isdn_traceflags} >/dev/null 2>&1 &
	fi
	echo '.'
}

load_rc_config $name
run_rc_command "$1"
