# etc/rc.d/isdnd
#!/bin/sh
#
# $NetBSD: isdnd,v 1.8 2002/03/22 04:33:58 thorpej Exp $
#

# PROVIDE: isdnd
# REQUIRE: NETWORKING syslogd mountcritremote
# BEFORE:  SERVERS
#
#	Note that this means that syslogd will not be listening on
#	any isdn addresses.  This is considered a feature.
#

. /etc/rc.subr

name="isdnd"
rcvar=$name
required_files="/etc/isdn/${name}.rc"

command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config $name
run_rc_command "$1"
