# etc/rc.d/kdc
#!/bin/sh
#
# $NetBSD: kdc,v 1.2 2000/09/19 13:04:38 lukem Exp $
#

# PROVIDE: kdc
# REQUIRE: NETWORK

. /etc/rc.subr

name="kdc"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/krb5.conf"

load_rc_config $name
run_rc_command "$1"
