# etc/rc.d/kdc
#!/bin/sh
#
# $NetBSD: kdc,v 1.5 2002/03/22 04:33:58 thorpej Exp $
# $FreeBSD$
#

# PROVIDE: kdc
# REQUIRE: NETWORKING
# BEFORE:  SERVERS

. /etc/rc.subr

name="kdc"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/krb5.conf"

load_rc_config $name
run_rc_command "$1"
