# etc/rc.d/kerberos
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/kerberos,v 1.8.2.1.4.1 2010/06/14 02:09:06 kensmith Exp $
#

# PROVIDE: kerberos
# REQUIRE: NETWORKING

. /etc/rc.subr

name="kerberos5"
rcvar="kerberos5_server_enable"

load_rc_config $name
command="${kerberos5_server}"
kerberos5_flags="${kerberos5_server_flags}"
run_rc_command "$1"
