# etc/rc.d/kerberos
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/kerberos,v 1.3.6.1 2004/10/10 09:50:53 mtm Exp $
#

# PROVIDE: kerberos
# REQUIRE: NETWORKING

. /etc/rc.subr

name="kerberos5"
load_rc_config $name
rcvar="kerberos5_server_enable"
command="${kerberos5_server}"
command_args="&"

run_rc_command "$1"
