# etc/rc.d/kerberos
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/kerberos,v 1.6.2.2.6.1 2010/02/10 00:26:20 kensmith Exp $
#

# PROVIDE: kerberos
# REQUIRE: NETWORKING

. /etc/rc.subr

name="kerberos5"
rcvar="kerberos5_server_enable"

load_rc_config $name
command="${kerberos5_server}"
kerberos5_flags="${kerberos5_server_flags}"
run_rc_command "$1"
