# etc/rc.d/kerberos
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/kerberos,v 1.3.6.2 2005/09/25 19:06:24 rwatson Exp $
#

# PROVIDE: kerberos
# REQUIRE: NETWORKING

. /etc/rc.subr

name="kerberos5"
load_rc_config $name
rcvar="kerberos5_server_enable"
command="${kerberos5_server}"
command_args="${kerberos5_server_flags} &"

run_rc_command "$1"
