# etc/rc.d/keyserv
#!/bin/sh
#
# $FreeBSD$
#

# Start keyserv if we are running Secure RPC
#

# PROVIDE: keyserv
# REQUIRE: ypbind
# BEFORE: DAEMON
# KEYWORD: FreeBSD

. /etc/rc.subr

name="keyserv"
rcvar=`set_rcvar`
command="/usr/sbin/${name}"
start_precmd="keyserv_prestart"

keyserv_prestart()
{
	if ! checkyesno rpcbind_enable  && \
		! /etc/rc.d/rpcbind forcestatus 1>/dev/null 2>&1
	then
		force_depend rpcbind || return 1
	fi

	return 0
}

load_rc_config $name
run_rc_command "$1"
