# etc/rc.d/keyserv
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/keyserv,v 1.6.2.3.2.1 2012/03/03 06:15:13 kensmith Exp $
#

# Start keyserv if we are running Secure RPC

# PROVIDE: keyserv
# REQUIRE: ypset
# BEFORE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="keyserv"
rcvar="keyserv_enable"
command="/usr/sbin/${name}"
start_precmd="keyserv_prestart"

keyserv_prestart()
{
	force_depend rpcbind || return 1
}

load_rc_config $name
run_rc_command "$1"
