# etc/rc.d/ldconfig
#!/bin/sh
#
# $NetBSD: ldconfig,v 1.4 2001/04/25 13:00:17 lukem Exp $
#

# PROVIDE: ldconfig
# REQUIRE: mountall
# BEFORE:  DAEMON

. /etc/rc.subr

name="ldconfig"
ldconfig_command="/sbin/ldconfig"
start_cmd="ldconfig_start"
stop_cmd=":"

ldconfig_start()
{
	if [ -f ${ldconfig_command} ]; then
		echo "Creating a.out runtime link editor directory cache."
		${ldconfig_command}
	fi
}

load_rc_config $name
run_rc_command "$1"
