# etc/rc.d/lkm2
#!/bin/sh
#
# $NetBSD: lkm2,v 1.4 2000/09/19 13:04:38 lukem Exp $
# $FreeBSD$
#

# PROVIDE: beforemountlkm
# REQUIRE: syslogd

. /etc/rc.subr

name="lkm2"
rcvar="lkm"
start_cmd="lkm2_start"
stop_cmd="lkm2_stop"

#		load kernel modules specified in /etc/lkm.conf
#
lkm2_start()
{
	if [ -r /etc/rc.lkm ]; then
		lkmstage=BEFOREMOUNT
		set start ; . /etc/rc.lkm
	else
		warn "/etc/rc.lkm not found; LKMs not loaded."
	fi
}

lkm2_stop()
{
	if [ -r /etc/rc.lkm ]; then
		lkmstage=BEFOREMOUNT
		set stop ; . /etc/rc.lkm
	fi
}

load_rc_config lkm
run_rc_command "$1"
