# etc/rc.d/lomac
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/lomac,v 1.2 2002/11/26 17:51:03 rwatson Exp $
#

# PROVIDE: lomac
# REQUIRE: mountcritremote
# BEFORE: syslogd
# KEYWORD: FreeBSD

. /etc/rc.subr

name="lomac"
rcvar=`set_rcvar`
start_cmd="kldload mac_lomac >/dev/null 2>&1"
stop_cmd=":"

load_rc_config $name
run_rc_command "$1"
