# etc/rc.d/lpd
#!/bin/sh
#
# $NetBSD: lpd,v 1.4 2000/09/19 13:04:38 lukem Exp $
#

# PROVIDE: lpd
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="lpd"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/printcap"

load_rc_config $name
run_rc_command "$1"
