# etc/rc.d/lpd
#!/bin/sh
#
# $NetBSD: lpd,v 1.3 2000/05/13 08:45:07 lukem Exp $
#

# PROVIDE: lpd
# REQUIRE: DAEMON

. /etc/rc.subr

name="lpd"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/printcap"

load_rc_config $name
run_rc_command "$1"
