# etc/rc.d/mixerctl
#!/bin/sh
#
# $NetBSD: mixerctl,v 1.2 2002/06/02 19:04:10 jmcneill Exp $
#

# PROVIDE: mixerctl
# REQUIRE: mountcritremote

. /etc/rc.subr

name="mixerctl"
start_cmd="mixerctl_start"
stop_cmd=":"

mixerctl_start()
{
	if [ -r /etc/mixerctl.conf ]; then
		echo "Setting mixerctl variables..."
		while read setting; do
			mixerctl -n -w $setting
		done < /etc/mixerctl.conf
	fi
}

load_rc_config $name
run_rc_command "$1"
