# etc/rc.d/mopd
#!/bin/sh
#
# $NetBSD: mopd,v 1.4 2000/09/19 13:04:38 lukem Exp $
#

# PROVIDE: mopd
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="mopd"
rcvar=$name
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
