# etc/rc.d/mountcritremote
#!/bin/sh
#
# $NetBSD: mountcritremote,v 1.2 2000/05/13 08:45:07 lukem Exp $
#

# PROVIDE: mountcritremote
# REQUIRE: NETWORK root mountcritlocal

. /etc/rc.subr

#	mount critical remote filesystems
#	(as specified in $critical_filesystems)
#
name="mountcritremote"
start_cmd="mount_critical_filesystems remote"
stop_cmd=":"

load_rc_config $name
run_rc_command "$1"
