# etc/rc.d/mountcritremote
#!/bin/sh
#
# $NetBSD: mountcritremote,v 1.6 2002/03/27 08:53:42 lukem Exp $
#

# PROVIDE: mountcritremote
# REQUIRE: NETWORKING root mountcritlocal

. /etc/rc.subr

name="mountcritremote"
start_cmd="mountcritremote_start"
stop_cmd=":"

mountcritremote_start()
{
	#	Mount critical filesystems that may be `remote'.
	#	(as specified in $critical_filesystems_remote)
	#	This usually includes /usr.
	#
	mount_critical_filesystems remote

}

load_rc_config $name
run_rc_command "$1"
