# etc/rc.d/moused
#!/bin/sh
#
# $NetBSD:  $
#

# PROVIDE: moused
# REQUIRE: DAEMON

. /etc/rc.subr

name="moused"
rcvar=$name
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
