# etc/rc.d/moused
#!/bin/sh
#
# $NetBSD: moused,v 1.1 2001/10/29 23:25:01 augustss Exp $
# $FreeBSD$
#

# PROVIDE: moused
# REQUIRE: DAEMON
# KEYWORD: FreeBSD NetBSD

. /etc/rc.subr

name=moused
rcvar=`set_rcvar`

case ${OSTYPE} in
FreeBSD)
	start_cmd="moused_start"
	;;
NetBSD)
	command="/usr/sbin/${name}"
	;;
esac

moused_start()
{
	echo -n 'Starting moused:'
	/usr/sbin/moused ${moused_flags} -p ${moused_port} -t ${moused_type}

	_mousechar_arg=
	case ${mousechar_start} in
	[Nn][Oo] | '')
		;;
	*)
		echo -n ' mousechar_start'
		_mousechar_arg="-M ${mousechar_start}"
		;;
	esac

	for ttyv in /dev/ttyv* ; do
		vidcontrol < ${ttyv} ${_mousechar_arg} -m on
	done
	echo '.'
}

load_rc_config $name
run_rc_command "$1"
