# etc/rc.d/mrouted
#!/bin/sh
#
# $NetBSD: mrouted,v 1.4 2000/05/13 08:45:07 lukem Exp $
#

# PROVIDE: mrouted
# REQUIRE: DAEMON

. /etc/rc.subr

name="mrouted"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
