# etc/rc.d/mrouted
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/mrouted,v 1.10.2.1.6.1 2010/02/10 00:26:20 kensmith Exp $
#

# PROVIDE: mrouted
# REQUIRE: netif routing cleanvar
# KEYWORD: nojail

. /etc/rc.subr

name="mrouted"
rcvar=`set_rcvar`
command="/usr/local/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
