# etc/rc.d/mroute6d
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/mroute6d,v 1.4.18.1 2008/11/25 02:59:29 kensmith Exp $
#

# PROVIDE: mroute6d
# REQUIRE: network_ipv6
# KEYWORD: nojail

. /etc/rc.subr

name="mroute6d"
rcvar=`set_rcvar`
command="/usr/local/sbin/pim6dd"

load_rc_config $name
run_rc_command "$1"
