# etc/rc.d/mrouted
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/mrouted,v 1.11.2.2.2.1 2012/03/03 06:15:13 kensmith Exp $
#

# PROVIDE: mrouted
# REQUIRE: netif routing cleanvar
# KEYWORD: nojail

. /etc/rc.subr

name="mrouted"
rcvar="mrouted_enable"
command="/usr/local/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/${name}.conf"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
