# etc/rc.d/mroute6d
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: mroute6d
# REQUIRE: network_ipv6
# KEYWORD: nojail

. /etc/rc.subr

name="mroute6d"
rcvar=`set_rcvar`
command="/usr/local/sbin/pim6dd"

load_rc_config $name
run_rc_command "$1"
