# etc/rc.d/natd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: natd
# KEYWORD: nostart nojail

. /etc/rc.subr
. /etc/network.subr

name="natd"
rcvar=`set_rcvar`
command="/sbin/${name}"
start_precmd="natd_precmd"
start_cmd="natd_start"

natd_precmd()
{
	if ! ${SYSCTL} net.inet.divert > /dev/null 2>&1; then
		if ! kldload ipdivert; then
			warn unable to load IPDIVERT module.
			return 1
		fi
	fi

	return 0
}

natd_start()
{
	dhcp_list="`list_net_interfaces dhcp`"
	for ifn in ${dhcp_list}; do
		case ${natd_interface} in
		${ifn})
			natd_flags="$natd_flags -dynamic"
			;;
		*)
			;;
		esac
	done
	if [ -n "${natd_interface}" ]; then
		if echo ${natd_interface} | \
		grep -q -E '^[0-9]+(\.[0-9]+){0,3}$'; then
			natd_flags="$natd_flags -a ${natd_interface}"
		else
			natd_flags="$natd_flags -n ${natd_interface}"
		fi
	fi
	echo -n ' natd'
	${natd_program:-/sbin/natd} ${natd_flags} ${natd_ifarg}
}

load_rc_config $name
run_rc_command "$1"
