# etc/rc.d/newsyslog
#!/bin/sh
#
# $NetBSD: newsyslog,v 1.5 2002/03/24 15:51:26 lukem Exp $
# $FreeBSD$
#

# PROVIDE: newsyslog
# REQUIRE: mountcritremote sysdb
# BEFORE:  syslogd SERVERS

. /etc/rc.subr

name="newsyslog"
rcvar=$name
required_files="/etc/newsyslog.conf"
start_cmd="newsyslog_start"
stop_cmd=":"

newsyslog_start()
{
	echo -n "Trimming log files:"
	/usr/bin/newsyslog -s $rc_flags
	echo " done."
}

load_rc_config $name
run_rc_command "$1"
