# etc/rc.d/nfslocking
#!/bin/sh
#
# $NetBSD: nfslocking,v 1.5 2002/03/22 04:33:59 thorpej Exp $
#

# PROVIDE: nfslocking
# REQUIRE: nfsd
# BEFORE:  DAEMON

. /etc/rc.subr

arg=$1

start_precmd='checkyesno nfs_server || checkyesno nfs_client ||
		[ -n "$rc_force" ]'
stop_precmd=$start_precmd
status_precmd=$start_precmd

name="statd"
rcvar=$name
command="/usr/sbin/rpc.${name}"
load_rc_config $name
run_rc_command "$arg"

name="lockd"
rcvar=$name
command="/usr/sbin/rpc.${name}"
load_rc_config $name
run_rc_command "$arg"
