# etc/rc.d/nfsserver
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/nfsserver,v 1.5 2006/12/31 10:37:18 yar Exp $
#

# PROVIDE: nfsserver
# REQUIRE: NETWORKING mountcritremote
# KEYWORD: nojail

. /etc/rc.subr

name="nfsserver"
rcvar="nfs_server_enable"
start_cmd=":"
stop_cmd=":"
required_modules="nfsserver"

load_rc_config $name
run_rc_command "$1"
