# etc/rc.d/nfsserver
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: nfsserver
# REQUIRE: NETWORKING mountcritremote
# KEYWORD: FreeBSD

. /etc/rc.subr

name="nfsserver"
rcvar="nfs_server_enable"
start_cmd="nfsserver_start"
stop_cmd=":"

# Load nfs modules if they were not compiled into the kernel
nfsserver_start()
{
	if ! sysctl vfs.nfsrv >/dev/null 2>&1; then
		if ! kldload nfsserver; then
			warn 'Could not load NFS server module'
			return 1
		fi
	fi
	return 0
}

load_rc_config $name
run_rc_command "$1"
