# etc/rc.d/nfsserver
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/nfsserver,v 1.5.10.1.6.1 2010/12/21 17:09:25 kensmith Exp $
#

# PROVIDE: nfsserver
# REQUIRE: NETWORKING mountcritremote
# KEYWORD: nojail

. /etc/rc.subr

name="nfsserver"
rcvar="nfs_server_enable"
start_cmd=":"
stop_cmd=":"
required_modules="nfsserver"

load_rc_config $name
run_rc_command "$1"
