# etc/rc.d/nfsuserd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/nfsuserd,v 1.1.2.2.2.1 2012/03/03 06:15:13 kensmith Exp $
#

# PROVIDE: nfsuserd
# REQUIRE: NETWORKING
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="nfsuserd"
rcvar="nfsuserd_enable"
command="/usr/sbin/${name}"
sig_stop="USR1"

load_rc_config $name

run_rc_command "$1"
