# etc/rc.d/nfsuserd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/nfsuserd,v 1.1.2.1.6.1 2010/12/21 17:09:25 kensmith Exp $
#

# PROVIDE: nfsuserd
# REQUIRE: NETWORKING
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="nfsuserd"
rcvar=`set_rcvar`
command="/usr/sbin/${name}"
sig_stop="USR1"

load_rc_config $name

run_rc_command "$1"
