# etc/rc.d/nscd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: cached
# REQUIRE: DAEMON
# BEFORE: LOGIN

#
# Add the following lines to /etc/rc.conf to enable cached:
#
# cached_enable="YES"
#
# See cached(8) for flags
#

. /etc/rc.subr

name=cached
rcvar=`set_rcvar`

command=/usr/sbin/cached

cached_enable=${cached_enable:-"NO"}
cached_pidfile=${cached_pidfile:-"/var/run/cached.pid"}
cached_flags=${cached_flags:-""}

load_rc_config $name
run_rc_command "$1"
