# etc/rc.d/nscd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: nscd
# REQUIRE: DAEMON
# BEFORE: LOGIN

#
# Add the following lines to /etc/rc.conf to enable cached:
#
# nscd_enable="YES"
#
# See nscd(8) for flags
#

. /etc/rc.subr

name=nscd
rcvar=`set_rcvar`

command=/usr/sbin/nscd
extra_commands="flush"
flush_cmd="${command} -I all"

nscd_enable=${nscd_enable:-"NO"}
nscd_pidfile=${nscd_pidfile:-"/var/run/nscd.pid"}
nscd_flags=${nscd_flags:-""}

load_rc_config $name
run_rc_command "$1"
