# etc/rc.d/othermta
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/othermta,v 1.3.24.1 2010/02/10 00:26:20 kensmith Exp $
#

# PROVIDE: mail
# REQUIRE: LOGIN

# XXX - TEMPORARY SCRIPT UNTIL YOU WRITE YOUR OWN REPLACEMENT.
#
. /etc/rc.subr

load_rc_config 'XXX'

if [ -n "${mta_start_script}" ]; then
	[ "${mta_start_script}" != "/etc/rc.sendmail" ] && \
	    sh ${mta_start_script} "$1"
fi
