# etc/rc.d/pfsync
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: pfsync
# REQUIRE: root mountcritlocal netif
# KEYWORD: nojail

. /etc/rc.subr

name="pfsync"
rcvar=`set_rcvar`
start_precmd="pfsync_prestart"
start_cmd="pfsync_start"
stop_cmd="pfsync_stop"
required_modules="pf"

pfsync_prestart()
{
	# XXX Currently pfsync cannot be a module as it must register
	# a network protocol in a static kernel table.
	if ! kldstat -q -m pfsync; then
		warn "pfsync(4) must be statically compiled in the kernel."
		return 1
	fi

	case "$pfsync_syncdev" in
	'')
		warn "pfsync_syncdev is not set."
		return 1
		;;
	esac
	return 0
}

pfsync_start()
{
	echo "Enabling pfsync."
	ifconfig pfsync0 syncdev $pfsync_syncdev $pfsync_ifconfig up
}

pfsync_stop()
{
	echo "Disabling pfsync."
	ifconfig pfsync0 -syncdev down
}

load_rc_config $name
run_rc_command "$1"
