# etc/rc.d/postrandom
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: postrandom
# REQUIRE: random FILESYSTEMS
# BEFORE: LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name="postrandom"
start_cmd="${name}_start"
stop_cmd=":"

# This will remove old entropy file.
# According to Bruce Schneier, this is strongly recommended in order
# to avoid using same ${entropy_file} across reboots.
# Reference: Chapter 10.6, Practical Cryptography, ISBN: 0-471-22357-3

postrandom_start()
{
	case ${entropy_dir} in
	[Nn][Oo])
		;;
	*)
		entropy_dir=${entropy_dir:-/var/db/entropy}
		if [ -d "${entropy_dir}" ]; then
			rm -f ${entropy_dir}/*
		fi
		boot_dir=${boot_dir:-/boot}
		if [ -d "${boot_dir}" ]; then
			rm -f ${boot_dir}/entropy
		fi
		rm -f /entropy
		rm -f /var/db/entropy-file
		;;
	esac
}

load_rc_config random
run_rc_command "$1"
