# etc/rc.d/powerd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: powerd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name="powerd"
rcvar=`set_rcvar`
command="/usr/sbin/powerd"

load_rc_config $name
command_args="${powerd_flags}"

run_rc_command "$1"
