# etc/rc.d/powerd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/powerd,v 1.2.4.1 2005/07/21 09:35:32 bruno Exp $
#

# PROVIDE: powerd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name="powerd"
rcvar=`set_rcvar`
command="/usr/sbin/powerd"

load_rc_config $name
run_rc_command "$1"
