# etc/rc.d/quota
#!/bin/sh
#
# $NetBSD: quota,v 1.2 2000/05/13 08:45:08 lukem Exp $
#

# PROVIDE: quota
# REQUIRE: mountall
# BEFORE:  DAEMON

. /etc/rc.subr

name="quota"
start_cmd="quota_start"
stop_cmd="/usr/sbin/quotaoff -a"

quota_start()
{
	echo -n "Checking quotas:"
	/usr/sbin/quotacheck -a
	echo " done."
	/usr/sbin/quotaon -a
}

load_rc_config $name
run_rc_command "$1"
