# etc/rc.d/racoon
#!/bin/sh
#
# $NetBSD: racoon,v 1.3 2002/03/22 04:33:59 thorpej Exp $
# $FreeBSD$
#

# PROVIDE: ike
# REQUIRE: isdnd kdc ppp
# BEFORE:  SERVERS

. /etc/rc.subr

name="racoon"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/racoon/racoon.conf"

load_rc_config $name
run_rc_command "$1"
