# etc/rc.d/racoon
#!/bin/sh
#
# $NetBSD: racoon,v 1.2 2001/04/26 03:55:16 lukem Exp $
#

# PROVIDE: ike
# REQUIRE: isdnd kdc ppp
# BEFORE:  SERVERS

. /etc/rc.subr

name="racoon"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/racoon/racoon.conf"

load_rc_config $name
run_rc_command "$1"
