# etc/rc.d/rarpd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/rarpd,v 1.6.14.1 2008/01/28 07:58:31 dougb Exp $
#

# PROVIDE: rarpd
# REQUIRE: DAEMON cleanvar
# BEFORE:  LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name="rarpd"
rcvar=`set_rcvar`
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/ethers"

load_rc_config $name
run_rc_command "$1"
