# etc/rc.d/rarpd
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/rarpd,v 1.7.2.2.2.1 2012/03/03 06:15:13 kensmith Exp $
#

# PROVIDE: rarpd
# REQUIRE: DAEMON cleanvar
# BEFORE:  LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name="rarpd"
rcvar="rarpd_enable"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/etc/ethers"

load_rc_config $name
run_rc_command "$1"
