# etc/rc.d/route6d
#!/bin/sh
#
# $NetBSD: route6d,v 1.3 2000/05/13 08:45:09 lukem Exp $
#

# PROVIDE: route6d
# REQUIRE: DAEMON

. /etc/rc.subr

name="route6d"
rcvar=$name
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
