# etc/rc.d/route6d
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: route6d
# REQUIRE: network_ipv6
# KEYWORD: nojail

. /etc/rc.subr

name="route6d"

# XXX - Executable may be in a different location. The $name variable
#       is different from the variable in rc.conf(5) so the
#       subroutines in rc.subr won't catch it. In this case, it
#	is also needed by the eval statement in the FreeBSD conditional.
#
load_rc_config $name

rcvar="ipv6_router_enable"
command="${ipv6_router:-/usr/sbin/${name}}"
eval ${name}_flags=\"${ipv6_router_flags}\"

run_rc_command "$1"
