# etc/rc.d/rpcbind
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/rpcbind,v 1.12.10.2.4.1 2010/02/10 00:26:20 kensmith Exp $
#

# PROVIDE: rpcbind
# REQUIRE: NETWORKING ntpdate syslogd named
# KEYWORD: shutdown

. /etc/rc.subr

name="rpcbind"
rcvar=`set_rcvar`
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
