# etc/rc.d/rpcbind
#!/bin/sh
#
# $NetBSD: rpcbind,v 1.5 2000/09/19 13:04:39 lukem Exp $
#

# PROVIDE: rpcbind
# REQUIRE: NETWORKING ntpdate syslogd named ppp

. /etc/rc.subr

name="rpcbind"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config $name
run_rc_command "$1"
