# etc/rc.d/rpcbind
#!/bin/sh
#
# $NetBSD: rpcbind,v 1.6 2002/01/31 01:26:06 lukem Exp $
# $FreeBSD$
#

# PROVIDE: rpcbind
# REQUIRE: NETWORKING ntpdate syslogd named ppp
# KEYWORD: FreeBSD NetBSD

. /etc/rc.subr

name="rpcbind"

# XXX - Executable may be in a different location. The $name variable
#	is different from the variable in rc.conf(5) so the
#	subroutines in rc.subr won't catch it.
#
load_rc_config $name
command="${portmap_program:-/usr/sbin/${name}}"

case `${CMD_OSTYPE}` in
FreeBSD)
	pidfile=
	rcvar="portmap_enable"
	eval ${name}_flags=\"${portmap_flags}\"
	;;
NetBSD)
	rcvar=$name
	command="/usr/sbin/${name}"
	pidfile="/var/run/${name}.pid"
	;;
esac

run_rc_command "$1"
