# etc/rc.d/rpcbind
#!/bin/sh
#
# $NetBSD: rpcbind,v 1.6 2002/01/31 01:26:06 lukem Exp $
# $FreeBSD$
#

# PROVIDE: rpcbind
# REQUIRE: NETWORKING ntpdate syslogd named ppp
# KEYWORD: FreeBSD NetBSD

. /etc/rc.subr

name="rpcbind"
rcvar=`set_rcvar`
command="/usr/sbin/${name}"

case ${OSTYPE} in
NetBSD)
	pidfile="/var/run/${name}.pid"
	;;
esac

load_rc_config $name
run_rc_command "$1"
