# etc/rc.d/rtadvd
#!/bin/sh
#
# $NetBSD: rtadvd,v 1.3 2000/05/13 08:45:09 lukem Exp $
#

# PROVIDE: rtadvd
# REQUIRE: DAEMON

. /etc/rc.subr

name="rtadvd"
rcvar=$name
command="/usr/sbin/${name}"
start_precmd="rtadvd_precmd"

rtadvd_precmd()
{
	if [ "$ip6mode" != "router" ]; then
		warn \
	"${name} cannot be used on IPv6 host, only on an IPv6 router."
		return 1
	fi
}

load_rc_config $name
run_rc_command "$1"
