# etc/rc.d/rwho
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/rwho,v 1.4.2.1.2.1 2008/10/02 02:57:24 kensmith Exp $
#

# PROVIDE: rwho
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="rwhod"
rcvar="`set_rcvar`"
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
