# etc/rc.d/rwho
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/rwho,v 1.4.10.2.6.1 2010/12/21 17:10:29 kensmith Exp $
#

# PROVIDE: rwho
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="rwhod"
rcvar="`set_rcvar`"
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
