# etc/rc.d/rwho
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/rwho,v 1.4.14.1 2008/01/28 07:58:31 dougb Exp $
#

# PROVIDE: rwho
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="rwhod"
rcvar="`set_rcvar`"
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
