# etc/rc.d/screenblank
#!/bin/sh
#
# $NetBSD: screenblank,v 1.4 2000/09/19 13:04:39 lukem Exp $
#

# PROVIDE: screenblank
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="screenblank"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config $name
run_rc_command "$1"
