# etc/rc.d/screenblank
#!/bin/sh
#
# $NetBSD: screenblank,v 1.5 2002/03/22 04:34:00 thorpej Exp $
# $FreeBSD$
#

# PROVIDE: screenblank
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="screenblank"
rcvar=$name
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

load_rc_config $name
run_rc_command "$1"
