# etc/rc.d/sdpd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: sdpd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name="sdpd"
control="${sdpd_control:-/var/run/sdp}"
group="${sdpd_group:-nobody}"
user="${sdpd_user:-nobody}"
command="/usr/sbin/${name}"
command_args="-c ${control} -g ${group} -u ${user}"
rcvar=`set_rcvar`

load_rc_config $name
run_rc_command "$1"
