# etc/rc.d/statd
#!/bin/sh
#
# $NetBSD: nfslocking,v 1.6 2002/03/24 15:52:41 lukem Exp $
# $FreeBSD$
#

# PROVIDE: nfslocking
# REQUIRE: nfsserver nfsclient nfsd
# BEFORE:  DAEMON
# KEYWORD: FreeBSD

. /etc/rc.subr

arg=$1

start_precmd='checkyesno nfs_server_enable || checkyesno nfs_client_enable ||
    [ -n "$rc_force" ]'
stop_precmd=$start_precmd
status_precmd=$start_precmd

name="statd"
rcvar=rpc_statd_enable
command="/usr/sbin/rpc.${name}"
load_rc_config $name
run_rc_command "$arg"

name="lockd"
rcvar=rpc_lockd_enable
command="/usr/sbin/rpc.${name}"
load_rc_config $name
run_rc_command "$arg"
