# etc/rc.d/swap1
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/swap1,v 1.13.2.1.4.1 2010/06/14 02:09:06 kensmith Exp $
#

# PROVIDE: localswap
# REQUIRE: disks
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="swap1"
start_cmd='swapon -aq'
stop_cmd=':'

load_rc_config swap
run_rc_command "$1"
