# etc/rc.d/swap1
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/swap1,v 1.9.10.1.6.1 2010/02/10 00:26:20 kensmith Exp $
#

# PROVIDE: localswap
# REQUIRE: disks
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="swap1"
start_cmd='swapon -a'
stop_cmd=':'

load_rc_config swap
run_rc_command "$1"
