# etc/rc.d/timed
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/timed,v 1.5.14.1 2008/01/28 07:58:32 dougb Exp $
#

# PROVIDE: timed
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name="timed"
rcvar=`set_rcvar`
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
