# etc/rc.d/timed
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/timed,v 1.5.2.1.2.1 2008/10/02 02:57:24 kensmith Exp $
#

# PROVIDE: timed
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: nojail

. /etc/rc.subr

name="timed"
rcvar=`set_rcvar`
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
