# etc/rc.d/timed
#!/bin/sh
#
# $FreeBSD: src/etc/rc.d/timed,v 1.5.10.2.4.1 2010/02/10 00:26:20 kensmith Exp $
#

# PROVIDE: timed
# REQUIRE: DAEMON
# BEFORE:  LOGIN
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="timed"
rcvar=`set_rcvar`
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
