# etc/rc.d/timed
#!/bin/sh
#
# $NetBSD: timed,v 1.4 2000/09/19 13:04:39 lukem Exp $
#

# PROVIDE: timed
# REQUIRE: DAEMON
# BEFORE:  LOGIN

. /etc/rc.subr

name="timed"
rcvar=$name
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
