# etc/rc.d/virecover
#!/bin/sh
#
# $NetBSD: virecover,v 1.3 2000/07/26 00:11:49 lukem Exp $
#

# PROVIDE: virecover
# REQUIRE: mountcritremote
#
# XXX: should require `mail'!

. /etc/rc.subr

name="virecover"
start_cmd="virecover_start" 
stop_cmd=":"

virecover_start()
{
	recover=`find /var/tmp/vi.recover -name 'recover.*' -type f -print`
	if [ -n "$recover" ]; then
		echo "Preserving editor files."
		for i in $recover; do
			sendmail -t < $i
		done
	fi
}

load_rc_config $name
run_rc_command "$1"
