# etc/rc.d/wpa_supplicant
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: wpa_supplicant
# REQUIRE: mountcritremote
# KEYWORD: nojail nostart

. /etc/rc.subr
. /etc/network.subr

name="wpa_supplicant"
rcvar=
command="/usr/sbin/${name}"
conf_file="/etc/wpa_supplicant.conf"

ifn="$2"
if [ -z "$ifn" ]; then
	return 1
fi

is_wired_interface()
{
	media=`ifconfig $1 2>/dev/null | while read line; do
		case "$line" in
		*media:?Ethernet*)
			echo Ethernet
			;;
		esac
	done`
	test "$media" = "Ethernet"
}

case ${ifn} in
ndis*)
	driver="ndis"
	;;
*)
	if is_wired_interface ${ifn} ; then
		driver="wired"
	else
		driver="bsd"
	fi
	;;
esac

load_rc_config $name

pidfile="/var/run/${name}/${ifn}.pid"
command_args="-B -q -i $ifn -c $conf_file -D $driver -P $pidfile"
required_files=$conf_file

run_rc_command "$1"
