# etc/rc.d/xdm
#!/bin/sh
#
# $NetBSD: xdm,v 1.5 2000/07/17 15:24:48 lukem Exp $
#

# PROVIDE: xdm
# REQUIRE: DAEMON LOGIN wscons
# KEYWORD: shutdown

. /etc/rc.subr

name="xdm"
rcvar=$name
command="/usr/X11R6/bin/${name}"
pidfile="/var/run/${name}.pid"
required_files="/usr/X11R6/lib/X11/xdm/xdm-config"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
