# etc/rc.d/xfs
#!/bin/sh
#
# $NetBSD: xfs,v 1.6 2002/03/22 04:34:01 thorpej Exp $
# $FreeBSD$
#

# PROVIDE: xfs
# REQUIRE: mountall cleartmp
# BEFORE:  LOGIN

. /etc/rc.subr

name="xfs"
rcvar=$name
command="/usr/X11R6/bin/${name}"
command_args="& sleep 2"
required_files="/usr/X11R6/lib/X11/fs/config"
extra_commands="reload"

load_rc_config $name
run_rc_command "$1"
