# etc/periodic/monthly/200.accounting
#!/bin/sh -
#
# $FreeBSD$
#

W=/var/log/wtmp
if [ -f $W ] ; then
	# Work out last month's name
	M=`date -v-1m +%h`

	mv -f $W $W.$M
	touch $W

	if [ -x /usr/sbin/ac ] ; then
		echo ""
		echo "Doing login accounting:"

		ac -p -w $W.$M | sort -nr +1
	fi

	rm -f $W.$M.gz
	gzip $W.$M
fi
