# etc/periodic/monthly/200.accounting
#!/bin/sh -
#
# $FreeBSD$
#

W=/var/log/wtmp
if [ -f ${W}.0 ] ; then
	if [ -x /usr/sbin/ac ] ; then
		echo ""
		echo "Doing login accounting:"

		ac -p -w ${W}.0 | sort -nr +1
	fi
fi
