# etc/periodic/monthly/200.accounting
#!/bin/sh -
#
# $FreeBSD: src/etc/periodic/monthly/200.accounting,v 1.4 1999/08/27 23:24:06 peter Exp $
#

W=/var/log/wtmp
if [ -f ${W}.0 ] ; then
	if [ -x /usr/sbin/ac ] ; then
		echo ""
		echo "Doing login accounting:"

		ac -p -w ${W}.0 | sort -nr +1
	fi
fi
