# etc/periodic/monthly/200.accounting
#!/bin/sh -
#
# $Id: 200.accounting,v 1.2 1998/02/09 11:36:43 danny Exp $
#

W=/var/log/wtmp
if [ -f ${W}.0 ] ; then
	if [ -x /usr/sbin/ac ] ; then
		echo ""
		echo "Doing login accounting:"

		ac -p -w ${W}.0 | sort -nr +1
	fi
fi
