# etc/periodic/monthly/200.accounting
#!/bin/sh -
#
# $Id: 200.accounting,v 1.1.1.1 1997/08/12 17:51:16 pst Exp $
#

if [ -x /usr/sbin/ac ] ; then
    echo ""
    echo "Doing login accounting:"

    ac -p | sort -nr +1
fi
