# etc/periodic/monthly/999.local
#!/bin/sh -
#
# $Id: 999.local,v 1.1.1.1 1997/08/16 17:04:01 pst Exp $
#
if [ -f /etc/monthly.local ]; then
    echo ""
    echo "Running monthly.local:"

    sh /etc/monthly.local
fi
