# etc/periodic/monthly/999.local
#!/bin/sh -
#
# $FreeBSD: src/etc/periodic/monthly/999.local,v 1.2 1999/08/27 23:24:06 peter Exp $
#
if [ -f /etc/monthly.local ]; then
    echo ""
    echo "Running monthly.local:"

    sh /etc/monthly.local
fi
