# etc/periodic/monthly/999.local
#!/bin/sh -
#
# $FreeBSD: src/etc/periodic/monthly/999.local,v 1.1.1.1.2.1 1999/08/29 14:19:48 peter Exp $
#
if [ -f /etc/monthly.local ]; then
    echo ""
    echo "Running monthly.local:"

    sh /etc/monthly.local
fi
