# etc/periodic/monthly/Makefile
# $FreeBSD$

FILES=	200.accounting \
	999.local

.include <bsd.prog.mk>
