# etc/periodic/monthly/Makefile
#	$Id: Makefile,v 1.1 1997/08/12 18:23:28 pst Exp $

BIN=	200.accounting \
	999.local

.include <bsd.prog.mk>
