# contrib/libpcap/rpcapd/rpcapd@.service
[Unit]
Description=Rpcap Per-Connection Server

[Service]
ExecStart=-/usr/local/sbin/rpcapd -i
StandardInput=socket
