# contrib/libpcap/rpcapd/rpcapd.socket
[Unit]
Description=Rpcap Socket for Per-Connection Servers

[Socket]
ListenStream=2002
Accept=yes

[Install]
WantedBy=sockets.target
