# etc/rc
#!/bin/sh
#
# $NetBSD: rc,v 1.152 2000/04/22 03:01:22 lukem Exp $
#
# rc.sh --
#	Run the scripts in /etc/rc.d with rcorder.

#	System startup script run by init on autoboot or after single-user.
#	Output and error are redirected to console by init, and the console
#	is the controlling terminal.

export HOME=/
export PATH=/sbin:/bin:/usr/sbin:/usr/bin

. /etc/rc.subr
. /etc/rc.conf

if ! checkyesno rc_configured; then
	echo "/etc/rc.conf is not configured.  Multiuser boot aborted."
	exit 1
fi

if [ "$1" = autoboot ]; then
	autoboot=yes
	_rc_fast_run=yes	# run_rc_command(): do fast booting
fi

stty status '^T'

#	Set shell to ignore SIGINT (2), but not children;
#	shell catches SIGQUIT (3) and returns to single user.
#
trap : 2
trap "echo 'Boot interrupted.'; exit 1" 3

files=`rcorder -s nostart /etc/rc.d/*`

for i in $files; do
	run_rc_script $i start
done

date
exit 0
