# etc/rc.local
#
# site-specific startup actions, daemons
#
#	$Id: rc.local,v 1.20 1996/08/17 07:15:38 peter Exp $
#

T=/tmp/_motd
rm -f $T
uname -v | sed -e 's,^\([^#]*\) #\(.* [1-2][0-9][0-9][0-9]\).*/\([^\]*\) $,\1 (\3) #\2,' > $T
awk '{if (NR == 1) {if ($1 == "FreeBSD") {next} else {print "\n"$0}} else {print}}' < /etc/motd >> $T
cp $T /etc/motd
chmod 644 /etc/motd
rm -f $T

echo -n 'starting local daemons:'

# put your local stuff here

echo '.'
