# etc/rc.local
#
# site-specific startup actions, daemons
#
#	$Id: rc.local,v 1.17 1995/03/30 00:10:31 ache Exp $
#

T=/tmp/_motd
rm -f $T
uname -v | sed -e 's/199[0-9]/&\
/' > $T
echo "" >> $T
sed '1,/^$/d' < /etc/motd >> $T
cp $T /etc/motd
chmod 644 /etc/motd
rm -f $T
_HOST=`hostname`

echo -n 'starting local daemons:'

# put your local stuff here

echo '.'
