# etc/rc.local
#
# site-specific startup actions, daemons
#
#	$Id: rc.local,v 1.16 1995/03/16 16:58:01 wollman Exp $
#

T=/tmp/_motd
rm -f $T
uname -rs > $T
echo "" >> $T
sed '1,/^$/d' < /etc/motd >> $T
cp $T /etc/motd
chmod 644 /etc/motd
rm -f $T
_HOST=`hostname`

echo -n 'starting local daemons:'

echo '.'
