# etc/rc.local
#
# site-specific startup actions, daemons
#
#	$Id$
#

T=/tmp/_motd
rm -f $T
uname -rs > $T
echo "" >> $T
sed '1,/^$/d' < /etc/motd >> $T
cp $T /etc/motd
chmod 644 /etc/motd
rm -f $T
_HOST=`hostname`

echo -n 'starting local daemons:'

# For loading fonts/keyboard example look in /usr/share/examples/syscons
# directory

echo '.'
