# etc/rc.local
#
# site-specific startup actions, daemons
#
#	$Id: rc.local,v 1.25 1997/06/25 15:15:13 ache Exp $
#
#	put your local stuff here or in /usr/local/etc/rc.d.  If you put stuff
#	here and want to make it programmer, bring in /etc/rc.conf and add
#	your custom variables to /etc/rc.conf.local
#
#	if [ -f /etc/rc.conf ]; then
#		. /etc/rc.conf
#	fi
#
#	... additional startup conditionals ...

