# 7/etc/rc.d/rpcbind
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: rpcbind
# REQUIRE: NETWORKING ntpdate syslogd named

. /etc/rc.subr

name="rpcbind"
rcvar=`set_rcvar`
command="/usr/sbin/${name}"

load_rc_config $name
run_rc_command "$1"
