# 6/etc/periodic/security/Makefile
# $FreeBSD$

FILES=	100.chksetuid \
	200.chkmounts \
	300.chkuid0 \
	400.passwdless \
	500.ipfwdenied \
	510.ipfdenied \
	520.pfdenied \
	550.ipfwlimit \
	600.ip6fwdenied \
	650.ip6fwlimit \
	700.kernelmsg \
	800.loginfail \
	900.tcpwrap \
	security.functions

.include <bsd.prog.mk>
