# 6/etc/rc.d/swap1
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: localswap
# REQUIRE: disks
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="swap1"
start_cmd='swapon -a'
stop_cmd=':'

load_rc_config swap
run_rc_command "$1"
