# release/rc.local
#!/bin/sh
# $FreeBSD$

# resolv.conf from DHCP ends up in here, so make sure the directory exists
mkdir /tmp/bsdinstall_etc

/usr/libexec/bsdinstall/startbsdinstall
