# etc/periodic/weekly/120.clean-kvmdb
#!/bin/sh -
#
# $FreeBSD: src/etc/periodic/weekly/120.clean-kvmdb,v 1.2.2.2 1999/08/29 14:19:49 peter Exp $
#

if [ -d /var/db ] ; then
    echo ""
    echo "Cleaning up kernel database files:"

    kernel=`sysctl -n kern.bootfile`
    kernel=kvm_`basename ${kernel}`.db

    find /var/db -name "kvm_*.db" -a ! -name ${kernel} -a \
	-atime +7 -delete
fi
