# etc/periodic/weekly/120.clean-kvmdb
#!/bin/sh -
#
# $Id: 120.clean-kvmdb,v 1.2 1998/01/01 09:04:47 bde Exp $
#

if [ -d /var/db ] ; then
    echo ""
    echo "Cleaning up kernel database files:"

    kernel=`sysctl -n kern.bootfile`
    kernel=kvm_`basename ${kernel}`.db

    find /var/db -name "kvm_*.db" -a ! -name ${kernel} -a \
	-atime +7 -delete
fi
