# etc/periodic/weekly/300.uucp
#!/bin/sh -
#
# $Id: 300.uucp,v 1.1.1.1 1997/08/16 17:04:02 pst Exp $
#
# This is really here for backwards compatibility, clean.weekly is not
# created by default anymore.

if [ -d /var/spool/uucp -a -f /usr/libexec/uucp/clean.weekly ]; then
    echo ""
    echo "Cleaning up UUCP:"

    echo /usr/libexec/uucp/clean.weekly | su daemon
fi
