# etc/periodic/weekly/310.locate
#!/bin/sh -
#
# $Id: 310.locate,v 1.1.1.1 1997/08/16 17:04:02 pst Exp $
#
locdb=/var/db/locate.database

if [ -x /usr/libexec/locate.updatedb -a -f $locdb ] ; then

    echo ""
    echo "Rebuilding locate database:"

    touch ${locdb}; chown nobody ${locdb}; chmod 644 ${locdb}

    echo /usr/libexec/locate.updatedb | nice -5 su -fm nobody
    chmod 444 ${locdb}
fi
