# etc/periodic/weekly/330.catman
#!/bin/sh -
#
# $Id: 330.catman,v 1.1.1.1 1997/08/12 17:51:16 pst Exp $
#

exit 0		# do not run by default

if [ -x /usr/libexec/catman.local -a -d /usr/share/man/cat1 ] ; then
    echo ""
    echo "Reformatting manual pages:"
    echo /usr/libexec/catman.local "${MANPATH}" | su -fm man
fi
