# etc/periodic/weekly/340.noid
#!/bin/sh -
#
# $FreeBSD$
#

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

case "$weekly_noid_enable" in
    [Yy][Ee][Ss])
	echo ""
	echo "Check for files with an unknown user or group:"

	find -H ${weekly_noid_dirs:-/} -fstype local \
	    \( -nogroup -o -nouser \) -print | sed 's/^/  /';;
esac
