# etc/periodic/weekly/340.noid
#!/bin/sh -
#
# $Id: 340.noid,v 1.1 1997/11/01 15:03:04 wosch Exp $

exit 0		# do not run by default

echo ""
echo "Check for files belongs to an unknown user or unknown group:"

# directories to be verified
: ${NOIDSEARCHPATHS="/"}

find -H $NOIDSEARCHPATHS -fstype local \( -nogroup -or -nouser \) -print

