# etc/periodic/weekly/999.local
#!/bin/sh -
#
# $Id: 999.local,v 1.2 1997/08/17 17:55:45 pst Exp $
#

if [ -f /etc/weekly.local ]; then
    echo ""
    echo "Running weekly.local:"

    sh /etc/weekly.local
fi
