# etc/periodic/weekly/999.local
#!/bin/sh -
#
# $FreeBSD: src/etc/periodic/weekly/999.local,v 1.2.2.1 1999/08/29 14:19:53 peter Exp $
#

if [ -f /etc/weekly.local ]; then
    echo ""
    echo "Running weekly.local:"

    sh /etc/weekly.local
fi
