# etc/periodic/weekly/999.local
#!/bin/sh -
#
# $FreeBSD$
#

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

for script in $weekly_local
do
    case "$script" in
	/*)
	    if [ -f "$script" ]
	    then
		echo ""
		echo "Running $script:"

		sh $script
	    fi;;
    esac
done
