# etc/periodic/weekly/999.local
#!/bin/sh -
#
# $Id: 999.local,v 1.1.1.1 1997/08/12 17:51:16 pst Exp $
#

if [ -f /etc/weekly.local ]; then
    echo ""
    echo "Running weekly.local:"
    sh /etc/weekly.local
fi
