# etc/periodic/weekly/999.local
#!/bin/sh -
#
# $FreeBSD$
#

if [ -f /etc/weekly.local ]; then
    echo ""
    echo "Running weekly.local:"

    sh /etc/weekly.local
fi
