# etc/periodic/weekly/Makefile
#	$Id: Makefile,v 1.1.1.1 1997/08/16 17:04:02 pst Exp $

BIN=	120.clean-kvmdb \
	300.uucp \
	310.locate \
	320.whatis \
	330.catman \
	999.local

.include <bsd.prog.mk>
