# etc/periodic/weekly/Makefile
# $FreeBSD: src/etc/periodic/weekly/Makefile,v 1.8 2002/07/18 12:33:01 ru Exp $

FILES=	120.clean-kvmdb \
	310.locate \
	320.whatis \
	330.catman \
	340.noid \
	400.status-pkg \
	999.local

.include <bsd.prog.mk>
