# etc/periodic/weekly/Makefile
#	$Id: Makefile,v 1.2 1997/08/12 23:11:24 adam Exp $

BIN=	100.clean-src \
	120.clean-kvmdb \
	300.uucp \
	310.locate \
	320.whatis \
	330.catman \
	999.local

.include <bsd.prog.mk>
