# etc/periodic/weekly/Makefile
# $FreeBSD$

FILES=	120.clean-kvmdb \
	310.locate \
	320.whatis \
	330.catman \
	340.noid \
	400.status-pkg \
	999.local

.include <bsd.prog.mk>
