# pf45/etc/rc.d/zvol
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: zvol
# REQUIRE: hostid
# KEYWORD: nojail

. /etc/rc.subr

name="zvol"
rcvar="zfs_enable"
start_cmd="zvol_start"
stop_cmd="zvol_stop"
required_modules="zfs"

zvol_start()
{
	zfs volinit
	# Enable swap on ZVOLs with property org.freebsd:swap=on.
	zfs list -H -o org.freebsd:swap,name -t volume | \
	while read state name; do
		case "${state}" in
		[oO][nN])
			swapon /dev/zvol/${name}
			;;
		esac
	done
}

zvol_stop()
{
	# Disable swap on ZVOLs with property org.freebsd:swap=on.
	zfs list -H -o org.freebsd:swap,name -t volume | \
	while read state name; do
		case "${state}" in
		[oO][nN])
			swapoff /dev/zvol/${name}
			;;
		esac
	done
	zfs volfini
}

load_rc_config $name
run_rc_command "$1"
